# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── mcp-server-fetch on port 8091 ───────────────────────────
; Server for fetching web pages
; https://pypi.org/project/mcp-server-fetch/
[program:mcp-fetch]
command=mcp-proxy --port=8091 --host=0.0.0.0 --allow-origin=* -- mcp-server-fetch
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-fetch.err

; ── mcp-server-time on port 8092 ────────────────────────────
; Server for date and time queries
; https://pypi.org/project/mcp-server-time/
[program:mcp-time]
command=mcp-proxy --port=8092 --host=0.0.0.0 --allow-origin=* -- mcp-server-time
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-time.err

; ── duckduckgo-mcp-server on port 8093 ──────────────────────
; Server for web search via DuckDuckGo
; Via mcp-proxy (stdio) – works around MCP SDK DNS rebinding issue (421 error)
; --fetch-backend curl enforces Chrome TLS fingerprint via curl_cffi
; https://pypi.org/project/duckduckgo-mcp-server/
[program:mcp-duckduckgo]
command=mcp-proxy --port=8093 --host=0.0.0.0 --allow-origin=* -- duckduckgo-mcp-server --fetch-backend curl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-duckduckgo.err

; ── mcp-file-edit on port 8094 ──────────────────────────────
; Flexible filesystem server with patch, code analysis and SSH features.
; No strict directory allowlist — use in controlled environments only!
; https://mcpservers.org/servers/patrickomatik/mcp-file-edit
; https://github.com/patrickomatik/mcp-file-edit
[program:mcp-file-edit]
command=mcp-proxy --port=8094 --host=0.0.0.0 --allow-origin=* -- python /opt/mcp-file-edit/server.py
directory=/data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-file-edit.err

; ── To add another server: ──────────────────────────────────
; [program:mcp-xyz]
; command=mcp-proxy --port=8096 --host=0.0.0.0 --allow-origin=* -- mcp-server-xyz
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stdout_logfile_maxbytes=0
; stderr_logfile=/tmp/mcp-xyz.err
