# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── mcp-server-fetch auf Port 8091 ──────────────────────────
[program:mcp-fetch]
command=mcp-proxy --port=8091 --host=0.0.0.0 --allow-origin=* -- mcp-server-fetch
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-fetch.err

; ── mcp-server-time auf Port 8092 ───────────────────────────
[program:mcp-time]
command=mcp-proxy --port=8092 --host=0.0.0.0 --allow-origin=* -- mcp-server-time
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-time.err

; ── duckduckgo-mcp-server auf Port 8093 ─────────────────────
; Via mcp-proxy (stdio) – umgeht das MCP-SDK DNS-Rebinding-Problem (421-Fehler)
; --fetch-backend curl erzwingt Chrome-TLS-Fingerprint via curl_cffi
[program:mcp-duckduckgo]
command=mcp-proxy --port=8093 --host=0.0.0.0 --allow-origin=* -- duckduckgo-mcp-server --fetch-backend curl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-duckduckgo.err

; ── mcp-filesystem auf Port 8094 ────────────────────────────
[program:mcp-filesystem]
command=mcp-proxy --port=8094 --host=0.0.0.0 --allow-origin=* -- mcp-filesystem --config /etc/mcp/filesystem.json
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-filesystem.err

; ── Weiteren Server hinzufügen: ─────────────────────────────
; [program:mcp-xyz]
; command=mcp-proxy --port=8095 --host=0.0.0.0 --allow-origin=* -- mcp-server-xyz
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stdout_logfile_maxbytes=0
; stderr_logfile=/tmp/mcp-xyz.err
