# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; ── mcp-server-fetch auf Port 8091 ──────────────────────────
[program:mcp-fetch]
command=mcp-proxy --port=8091 --host=0.0.0.0 --allow-origin=* mcp-server-fetch
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── mcp-server-time auf Port 8092 ───────────────────────────
[program:mcp-time]
command=mcp-proxy --port=8092 --host=0.0.0.0 --allow-origin=* mcp-server-time
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Weiteren Server hinzufügen: ─────────────────────────────
; [program:mcp-xyz]
; command=mcp-proxy --port=8093 --host=0.0.0.0 --allow-origin=* mcp-server-xyz
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stdout_logfile_maxbytes=0
; stderr_logfile=/dev/stderr
; stderr_logfile_maxbytes=0
