# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── mcp-server-fetch auf Port 8091 ──────────────────────────
; Server für Zugriffe auf Webseiten
; https://pypi.org/project/mcp-server-fetch/
[program:mcp-fetch]
command=mcp-proxy --port=8091 --host=0.0.0.0 --allow-origin=* -- mcp-server-fetch
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-fetch.err

; ── mcp-server-time auf Port 8092 ───────────────────────────
; Server für Datum und Uhrzeit
; https://pypi.org/project/mcp-server-time/
[program:mcp-time]
command=mcp-proxy --port=8092 --host=0.0.0.0 --allow-origin=* -- mcp-server-time
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-time.err

; ── duckduckgo-mcp-server auf Port 8093 ─────────────────────
; Server für Suche im Internet (Duck Duck Go)
; Via mcp-proxy (stdio) – umgeht das MCP-SDK DNS-Rebinding-Problem (421-Fehler)
; --fetch-backend curl erzwingt Chrome-TLS-Fingerprint via curl_cffi
; https://pypi.org/project/duckduckgo-mcp-server/
[program:mcp-duckduckgo]
command=mcp-proxy --port=8093 --host=0.0.0.0 --allow-origin=* -- duckduckgo-mcp-server --fetch-backend curl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-duckduckgo.err

; ── mcp-file-edit auf Port 8094 ─────────────────────────────
; Flexiblerer Filesystem-Server mit Patch-, Code-Analyse- und SSH-Features.
; Kein hartes Verzeichnis-Allowlist — nur für kontrollierte Umgebungen!
; https://mcpservers.org/servers/patrickomatik/mcp-file-edit
; https://github.com/patrickomatik/mcp-file-edit
[program:mcp-file-edit]
command=mcp-proxy --port=8094 --host=0.0.0.0 --allow-origin=* -- python /opt/mcp-file-edit/server.py
directory=/data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/mcp-file-edit.err

; ── Weiteren Server hinzufügen: ─────────────────────────────
; [program:mcp-xyz]
; command=mcp-proxy --port=8096 --host=0.0.0.0 --allow-origin=* -- mcp-server-xyz
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stdout_logfile_maxbytes=0
; stderr_logfile=/tmp/mcp-xyz.err
